# docker/supervisord.conf
# Supervisor 配置文件
# 管理: Code-Server + AionUi WebUI + VNC Server + Chrome CDP
# 版本: 4.0

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

# === Code-Server 进程配置 ===
[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth password --disable-telemetry
directory=/home/coder
user=coder
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/coder",PASSWORD="%(ENV_PASSWORD)s"

# === AionUi WebUI 进程配置 ===
[program:aionui-webui]
command=AionUi --webui --remote --no-sandbox
directory=/home/coder
user=coder
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/coder",DISPLAY="%(ENV_DISPLAY)s",AIONUI_PORT="%(ENV_AIONUI_PORT)s"

# === VNC Server 进程配置 ===
# AionUi 需要 X11 显示服务器，因此 VNC 必须自动启动
[program:vnc-server]
command=/vnc-entrypoint.sh
directory=/home/coder
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",DISPLAY="%(ENV_DISPLAY)s",VNC_PORT="%(ENV_VNC_PORT)s",NOVNC_PORT="%(ENV_NOVNC_PORT)s",VNC_GEOMETRY="%(ENV_VNC_GEOMETRY)s",VNC_DEPTH="%(ENV_VNC_DEPTH)s"

# === Chrome CDP 进程组配置 ===
[group:chrome-cdp]
programs=chrome-cdp-headless
priority=100

# === Chrome CDP (Headless 模式) ===
[program:chrome-cdp-headless]
command=/chrome-entrypoint.sh
directory=/home/coder
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",DISPLAY="%(ENV_DISPLAY)s",CDP_PORT="%(ENV_CDP_PORT)s",CHROME_PROFILE_DIR="%(ENV_CHROME_PROFILE_DIR)s",VNC_MODE="false",CHROME_BIN="%(ENV_CHROME_BIN)s"

# === Chrome CDP (VNC 模式) ===
[program:chrome-cdp-vnc]
command=/chrome-entrypoint.sh
directory=/home/coder
user=root
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",DISPLAY="%(ENV_DISPLAY)s",CDP_PORT="%(ENV_CDP_PORT)s",CHROME_PROFILE_DIR="%(ENV_CHROME_PROFILE_DIR)s",VNC_MODE="true",CHROME_BIN="%(ENV_CHROME_BIN)s"
